# containers/supervisor/datalab.conf
[supervisord]
nodaemon=true

[program:datalab]
autostart=true
autorestart=true
command=node /datalab/web/app.js
stdout_logfile=/var/log/datalab.stdout.log
redirect_stderr=true

