# containers/supervisor/gateway.conf
[supervisord]
nodaemon=true

[program:gateway]
autostart=true
autorestart=true
command=node /datalab/web/kernelproxy/app.js
stdout_logfile=/var/log/gateway.stdout.log
redirect_stderr=true

